# systemd/diva-embedded.service
[Unit]
Description=Diva Embedded Voice Assistant
After=network.target diva-audio.service
Requires=diva-audio.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/diva-embedded
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
