# systemd/diva-backup.service
[Unit]
Description=Diva Daily Encrypted Backup

[Service]
Type=oneshot
ExecStart=/opt/diva-watchdog/backup.sh
