# systemd/diva-wake.service
[Unit]
Description=Diva PROTO — Python wake word + audio
After=diva-node.service
Requires=diva-node.service

[Service]
Type=simple
WorkingDirectory=/opt/diva-current
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 /opt/diva-current/python/wakeword_server.py
Restart=always
RestartSec=5
Environment=AUDIO_INPUT_DEVICE=plughw:5
EnvironmentFile=/opt/diva-current/.env

[Install]
WantedBy=multi-user.target
