# systemd/diva-node.service
[Unit]
Description=Diva PROTO — Node.js API server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/diva-current
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/opt/diva-current/.env

[Install]
WantedBy=multi-user.target
